# backEnd/here_dev/sh_script/crontab

58 23 * * * $HOME/here/backEnd/here_dev/sh_script/errorLogReport.sh
59 23 * * * $HOME/here/backEnd/here_dev/sh_script/zipHereLog.sh

